# ServiceFiles/back-up-vimrc.service
[Unit]
Description=Full system upgrade

[Service]
Environment=SCRIPTS=/home/blake/scripts/autoscripts
#Environment=LOGS=/home/blake/Documents/logs
ExecStart=/bin/bash -c "${SCRIPTS}/back_up_vimrc"
